# keyboard-backlight/kbsv.service
[Unit]
Description=Keyboard backlight management

[Service]
ExecStart=/etc/kbsv/kbsv.py

[Install]
WantedBy=sysinit.target
